# lamp-controller.service
# To test:
#           journalctl -u lamp-controller.service -f
#
# To deploy:
#           sudo cp lamp-controller.service /etc/systemd/system/
#           sudo systemctl daemon-reload
#           sudo systemctl enable lamp-controller.service
#           sudo systemctl start lamp-controller.service

[Unit]
Description=Raspberry Pi Lamp and Sensor Controller
After=network.target
Wants=network.target

# If you want to be explicit about I2C:
# After=i2c-dev.service
# Wants=i2c-dev.service

[Service]
Type=simple
User=pi
Group=pi

# Set working directory so run.py finds lamp.config and logs locally
WorkingDirectory=/home/pi/lamp-controller

# Start the controller
ExecStart=/usr/bin/python3 /home/pi/lamp-controller/run.py

Restart=always
RestartSec=5

# Give it time for clean GPIO/sensor shutdown
TimeoutStopSec=20

Environment=PYTHONUNBUFFERED=1

# All print/logging goes to journal (plus your file logger inside run.py)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
